# scripts/com.kemmer.df-pms-mews-adapter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- DF-PMS-MEWS-Adapter macOS LaunchAgent [CRUX-MK]
     Welle-36 Foundation-DF (12. Foundation-DF im Kemmer-System).
     HeyLou-Mosaic-Layer fuer MEWS PMS-Integration.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kemmer.df-pms-mews-adapter</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/make/Projects/dark-factories/df-pms-mews-adapter/scripts/run-df-pms-mews-adapter-mac.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>7200</integer>

    <key>StandardOutPath</key>
    <string>/Users/make/Projects/dark-factories/df-pms-mews-adapter/runs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/make/Projects/dark-factories/df-pms-mews-adapter/runs/stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/make/Projects/dark-factories/df-pms-mews-adapter</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DF_PMS_MEWS_REAL_ENABLED</key>
        <string>false</string>
        <key>DF_PMS_MEWS_TENANT_ID</key>
        <string>hildesheim</string>
        <key>DF_PMS_MEWS_TIMEOUT_S</key>
        <string>30</string>
    </dict>
</dict>
</plist>
